# pi_cam_server/camera-server.service
# Pi Camera Server Service
[Unit]
Description=Pi Camera Server for Robot System
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/robot_camera
Environment=PYTHONPATH=/home/pi/robot_camera
ExecStart=/home/pi/robot_camera/venv/bin/python /home/pi/robot_camera/camera_server.py --daemon
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=camera-server

# Security
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/robot_camera

[Install]
WantedBy=multi-user.target
